# docs/mqtt-exporter-embedded.service
[Unit]
Description=MQTT Exporter with embedded broker for Meshtastic devices
After=network.target

[Service]
Type=simple
User=mqtt-exporter
Group=mqtt-exporter
WorkingDirectory=/opt/mqtt-exporter
ExecStart=/opt/mqtt-exporter/meshtastic-exporter-embedded --config /etc/mqtt-exporter/config.yaml
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/mqtt-exporter

[Install]
WantedBy=multi-user.target